# hour-timer.target
[Unit]
Description=hour-timer-target 
SopWhenUnneeded=yes
