# lab-time.timer
[Unit]
Description=Minute Timer
[Timer]
OnCalendar=*_*_* *:00:00
[Install]
WantedBy=basic.target